# rebootbench/experiments/tinyserver/rebootbench-tinyserver.service.tmpl
[Unit]
Description=rebootbench tinyserver SUT

[Service]
# Replace ABSOLUTE_PATH_TO_TINYSERVER before installing.
ExecStart=ABSOLUTE_PATH_TO_TINYSERVER
Restart=always
RestartSec=0
Environment=LISTEN_ADDR=:18080

[Install]
WantedBy=default.target
